# deploy/cloudpanel/supervisor-radiushub.conf
; Substitua __PROJECT_ROOT__, __SITE_USER__ e __PHP_BIN__ antes de instalar.
[program:radiushub-worker]
process_name=%(program_name)s_%(process_num)02d
command=__PHP_BIN__ __PROJECT_ROOT__/artisan queue:work --queue=webhooks,default --sleep=2 --tries=5 --timeout=120 --max-time=3600
directory=__PROJECT_ROOT__
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=__SITE_USER__
numprocs=1
redirect_stderr=true
stdout_logfile=__PROJECT_ROOT__/storage/logs/worker.log
stopwaitsecs=3600
